# etc/periodic/daily/backup-daily.ini
[program:backup-daily]
command=/usr/sbin/mysqltos3 daily
autostart=false
autorestart=false
stdout_events_enabled=true
stderr_events_enabled=true