# etc/periodic/daily/backup-daily
#!/bin/sh
supervisorctl start backup-daily