# etc/periodic/hourly/backup-hourly
#!/bin/sh
supervisorctl start backup-hourly