# etc/periodic/monthly/backup-monthly.ini
[program:backup-monthly]
command=/usr/sbin/mysqltos3 monthly
autostart=false
autorestart=false
stdout_events_enabled=true
stderr_events_enabled=true