# etc/cron.d/backup-monthly
30 3 1 * * root supervisorctl start backup-monthly