# etc/supervisor.d/backup-hourly.ini
[program:backup-hourly]
command=/usr/sbin/mysqltos3 hourly
autostart=false
autorestart=false
stdout_events_enabled=true
stderr_events_enabled=true