# etc/supervisor.d/cron.ini
[program:cron]
command=/usr/sbin/crond -f
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true