# etc/periodic/weekly/backup-weekly
#!/bin/sh
supervisorctl start backup-weekly