# etc/periodic/weekly/backup-weekly.ini
[program:backup-weekly]
command=/usr/sbin/mysqltos3 weekly
autostart=false
autorestart=false
stdout_events_enabled=true
stderr_events_enabled=true